# res/DEBIAN/prerm
#!/bin/bash

set -e

case $1 in
    remove|upgrade)
		INITSYS=$(ls -al /proc/1/exe | awk -F' ' '{print $NF}' | awk -F'/' '{print $NF}')
		rm -f /usr/bin/xldesk
		
		if [ "systemd" == "${INITSYS}" ]; then

			systemctl stop xldesk || true
			systemctl disable xldesk || true
			rm /etc/systemd/system/xldesk.service /usr/lib/systemd/system/xldesk.service || true
			
			# workaround temp dev build between 1.1.9 and 1.2.0
			serverUser=$(ps -ef | grep -E 'xldesk +--server' | grep -v 'sudo ' | awk '{print $1}' | head -1)
			if [ "$serverUser" != "" ] && [ "$serverUser" != "root" ]
			then
				systemctl --machine=${serverUser}@.host --user stop xldesk >/dev/null 2>&1 || true
			fi
			rm /usr/lib/systemd/user/xldesk.service >/dev/null 2>&1 || true
		fi
        ;;
esac

exit 0
